# deploy/systemd/hoyo-source-check.service
[Unit]
Description=HoYo FOMO Tracker Source Change Checker
After=network-online.target hoyo-api.service
Wants=network-online.target

[Service]
Type=oneshot
User=nasana
WorkingDirectory=/home/nasana/hoyo-tracker
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/nasana/hoyo-tracker/.venv/bin/python -m scripts.check_sources --all
